# sdcard/rootfs_overlay/etc/udev/rules.d/50-wifi.rules
KERNEL=="wlan0", RUN+="/sbin/hotplug.wlan"
