# rpi/rootfs_overlay/etc/udev/rules.d/50-dvb.rules
SUBSYSTEM=="dvb", RUN+="/sbin/hotplug.dvb"
