# rpi/rootfs_overlay/etc/udev/rules.d/50-storage.rules
ENV{DEVTYPE}=="partition", RUN+="/sbin/hotplug.storage"
