# rpi/rootfs_overlay/etc/udev/rules.d/50-storage.rules
ACTION=="add", KERNEL=="sd[a-z]*", RUN+="/sbin/hotplug.scsi"
