# rpi/rootfs_overlay/etc/udev/rules.d/50-storage.rules
KERNEL!="mmcblk0*", ENV{DEVTYPE}=="partition", RUN+="/sbin/hotplug.storage"
