# rpi/rootfs_overlay/etc/udev/rules.d/90-storage-hotplug.rules
ENV{ID_FS_USAGE}=="filesystem", KERNEL!="mmcblk*", KERNEL!="loop*", RUN+="/sbin/hotplug.storage"
